# padre-dashboard.service
[Unit]
Description=Padre Tracker Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cokle
WorkingDirectory=/home/cokle/padre-tracker/src
ExecStart=/home/cokle/padre-tracker/.venv/bin/gunicorn \
    --workers 4 \
    --worker-class gthread \
    --threads 8 \
    --timeout 60 \
    --bind 0.0.0.0:8080 \
    --access-logfile - \
    dashboard:app
Restart=always
RestartSec=10
# EnvironmentFile= ALWAYS overrides Environment= regardless of order — so the
# server's .env MUST NOT contain DB_PATH (or it would override the absolute
# path below and resolve relative to ./src, where no DB exists). We only want
# .env for HELIUS_API_KEY. Leading `-` = don't fail if missing.
EnvironmentFile=-/home/cokle/padre-tracker/.env
Environment=DASHBOARD_PORT=8080
Environment=DB_PATH=/home/cokle/padre-tracker/data/calls.db

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
