# padre-dashboard.service
[Unit]
Description=Padre Tracker Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cokle
WorkingDirectory=/home/cokle/padre-tracker/src
ExecStart=/home/cokle/padre-tracker/.venv/bin/gunicorn \
    --workers 4 \
    --worker-class gthread \
    --threads 8 \
    --timeout 60 \
    --bind 0.0.0.0:8080 \
    --access-logfile - \
    dashboard:app
Restart=always
RestartSec=10
Environment=DASHBOARD_PORT=8080
Environment=DB_PATH=/home/cokle/padre-tracker/data/calls.db
# Optional: .env provides HELIUS_API_KEY (for Helius DAS metadata fetching).
# Leading `-` = don't fail if missing — metadata worker is no-op without key.
EnvironmentFile=-/home/cokle/padre-tracker/.env

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
